# deploy/myapp.service
[Unit]
Description=MyApp Node Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/node /opt/myapp/index.js
WorkingDirectory=/opt/myapp
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
